# deploy/systemd/shaka-server.service
[Unit]
Description=Shaka Server Spring Boot Application
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu

EnvironmentFile=/etc/shaka/env
WorkingDirectory=/opt/shaka

ExecStart=/usr/bin/java \
    -Xmx320m -Xms256m \
    -XX:+UseG1GC \
    -XX:MaxGCPauseMillis=200 \
    -Djava.security.egd=file:/dev/./urandom \
    -jar /opt/shaka/current.jar

SuccessExitStatus=143
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/shaka /opt/shaka /tmp

[Install]
WantedBy=multi-user.target
